# templates/cron/systemd-timer.service
# Skeleton systemd service unit for a scheduled OS-level task.
# Pair with systemd-timer.timer.
# Install: drop both files into ~/.config/systemd/user/, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now <name>.timer

[Unit]
Description=<one-line description of what this does>
# Wants/After can chain dependencies if you need them, e.g.
# After=network-online.target
# Wants=network-online.target

[Service]
Type=oneshot
# Use absolute paths - systemd's PATH is minimal.
ExecStart=/path/to/your/script.sh
# Optional: load env vars from a file (rotated separately from this unit)
# EnvironmentFile=%h/.config/<name>/env
# Optional: cap runtime so a hung script doesn't pin a timer slot
TimeoutStartSec=10min
